# openwrt/target/linux/x86_64/base-files/etc/init.d/defconfig
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2010 Thinktube Inc.
# Copyright (C) 2009-2011 ETegro Technologies, PLC
#                         Sergey Matveev <stargrave@stargrave.org>
#

START=05

start() {
	. /lib/soekris.sh

	local board=astor2

	for f in $( ls /etc/defconfig/$board ); do
		if [ ! -e /etc/config/$f ]; then
			cp /etc/defconfig/$board/$f /etc/config/
		fi
	done
}
