# feeds/astor2/iscsi-scst-userspace/files/scst.init
#!/bin/sh /etc/rc.common
START=50

VAR_PR=/var/lib/scst/pr
SCSTD=/usr/sbin/iscsi-scstd

start()
{
	mkdir -p $VAR_PR
	$SCSTD
}

stop()
{
	true
}
