# feeds/astor2/astor2-luci-bwc-ipmi/files/astor2-luci-bwc-ipmi.init
#!/bin/sh /etc/rc.common

START=95
STOP=95

BWC=/usr/bin/luci-bwc-ipmi
BWC_DB=/tmp/lib/luci-bwc-ipmi

start() {
	mkdir -p $BWC_DB
	$BWC -d
}

stop() {
	killall ${BWC##*/}
}

restart() {
    stop
    start
}
