# feeds/luci/themes/astor2-san/root/etc/uci-defaults/luci-theme-astor2-san
#!/bin/sh
uci batch <<-EOF
	set luci.themes.astor2-san=/luci-static/astor2-san
        commit luci
EOF

