# feeds/astor2/astor2-lua-scst/files/astor2-lua-scst.init
#!/bin/sh /etc/rc.common
START=50

start()
{
	echo 'scst = require "astor2.scst";
	      scst.Daemon.apply()' | lua -
}

stop()
{
	true
}
