# feeds/luci/themes/astor2/root/etc/uci-defaults/luci-theme-astor2
#!/bin/sh
uci batch <<-EOF
	set luci.themes.astor2=/luci-static/astor2
	set luci.main.mediaurlbase=/luci-static/astor2
	set luci.main.lang=ru
        commit luci
EOF

