# build/cron.daily/astor2-build.sh
#!/bin/sh

cd /home/crontab/astor2/build
./build.sh
