# feeds/astor2/flashcache/files/flashcache.init
#!/bin/sh /etc/rc.common
START=11

start()
{
	insmod /lib/modules/flashcache.ko
}

stop()
{
	true
}
