# astor2-feed/astor2-einarc/files/einarc.init
#!/bin/sh /etc/rc.common
START=50

EINARC_STARTED_FLAG="/tmp/einarc-started"

start()
{
	einarc -t software -a 0 adapter restart
	touch $EINARC_STARTED_FLAG
}

stop()
{
	true
}

restart()
{
	stop
	start
}
